# udev/usr/lib/udev/rules.d/99-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", HOTPLUG=="1", RUN+="/usr/local/bin/udev-monitor"
